# sigma/pre_compiled/linuxbuiltin_savedsearches.conf
[default]
cron_schedule = */30 * * * *
dispatch.earliest_time = 0
dispatch.latest_time = now
enableSched = 0
schedule_window = auto
[Potential Suspicious BPF Activity - Linux]
description = Detects the presence of "bpf_probe_write_user" BPF helper-generated warning messages. Which could be a sign of suspicious eBPF activity on the system.
search = index=syslog _index_earliest=-1h@h "bpf_probe_write_user" | collect index=notable_events source="Potential Suspicious BPF Activity - Linux" marker="guid=0fadd880-6af3-4610-b1e5-008dc3a11b8a,tags=attack.persistence,tags=attack.defense-evasion,"
[Privileged User Has Been Created]
description = Detects the addition of a new user to a privileged group such as "root" or "sudo"
search = index=syslog _index_earliest=-1h@h "new user" "GID=0" OR "UID=0" OR "GID=10" OR "GID=27" | collect index=notable_events source="Privileged User Has Been Created" marker="guid=0ac15ec3-d24f-4246-aa2a-3077bb1cf90e,tags=attack.persistence,tags=attack.t1136.001,tags=attack.t1098,"
[Nimbuspwn Exploitation]
description = Detects exploitation of Nimbuspwn privilege escalation vulnerability (CVE-2022-29799 and CVE-2022-29800)
search = index=syslog _index_earliest=-1h@h "networkd-dispatcher" "Error handling notification for interface" "../../" | collect index=notable_events source="Nimbuspwn Exploitation" marker="guid=7ba05b43-adad-4c02-b5e9-c8c35cdf9fa8,tags=attack.privilege-escalation,tags=attack.t1068,"
[Shellshock Expression]
description = Detects shellshock expressions in log files
search = index=syslog _index_earliest=-1h@h "(){:;};" OR "() {:;};" OR "() { :;};" OR "() { :; };" | collect index=notable_events source="Shellshock Expression" marker="guid=c67e0c98-4d39-46ee-8f6b-437ebf6b950e,tags=attack.persistence,tags=attack.t1505.003,"
[Symlink Etc Passwd]
description = Detects suspicious command lines that look as if they would create symbolic links to /etc/passwd
search = index=syslog _index_earliest=-1h@h "ln -s -f /etc/passwd" OR "ln -s /etc/passwd" | collect index=notable_events source="Symlink Etc Passwd" marker="guid=c67fc22a-0be5-4b4f-aad5-2b32c4b69523,tags=attack.t1204.001,tags=attack.execution,"
[Suspicious Reverse Shell Command Line]
description = Detects suspicious shell commands or program code that may be executed or used in command line to establish a reverse shell
search = index=syslog _index_earliest=-1h@h "BEGIN {s = \"/inet/tcp/0/" OR "bash -i >& /dev/tcp/" OR "bash -i >& /dev/udp/" OR "sh -i >$ /dev/udp/" OR "sh -i >$ /dev/tcp/" OR "&& while read line 0<&5; do" OR "/bin/bash -c exec 5<>/dev/tcp/" OR "/bin/bash -c exec 5<>/dev/udp/" OR "nc -e /bin/sh " OR "/bin/sh | nc" OR "rm -f backpipe; mknod /tmp/backpipe p && nc " OR ";socket(S,PF_INET,SOCK_STREAM,getprotobyname(\"tcp\"));if(connect(S,sockaddr_in($p,inet_aton($i))))" OR ";STDIN->fdopen($c,r);$~->fdopen($c,w);system$_ while<>;" OR "/bin/sh -i <&3 >&3 2>&3" OR "uname -a; w; id; /bin/bash -i" OR "$sendbyte = ([text.encoding]::ASCII).GetBytes($sendback2); $stream.Write($sendbyte,0,$sendbyte.Length); $stream.Flush()};" OR ";os.dup2(s.fileno(),0);os.dup2(s.fileno(),1);os.dup2(s.fileno(),2);os.putenv('HISTFILE','/dev/null');" OR ".to_i;exec sprintf(\"/bin/sh -i <&%d >&%d 2>&%d\",f,f,f)" OR ";while(cmd=c.gets);IO.popen(cmd,\"r\"){|io|c.print" OR "socat exec:'bash -li',pty,stderr,setsid,sigint,sane tcp:" OR "rm -f /tmp/p; mknod /tmp/p p &&" OR " | /bin/bash | telnet " OR ",echo=0,raw tcp-listen:" OR "nc -lvvp " OR "xterm -display 1" | collect index=notable_events source="Suspicious Reverse Shell Command Line" marker="guid=738d9bcf-6999-4fdb-b4ac-3033037db8ab,tags=attack.execution,tags=attack.t1059.004,"
[Equation Group Indicators]
description = Detects suspicious shell commands used in various Equation Group scripts and tools
search = index=syslog _index_earliest=-1h@h "chown root*chmod 4777 " OR "cp /bin/sh .;chown" OR "chmod 4777 /tmp/.scsi/dev/bin/gsh" OR "chown root:root /tmp/.scsi/dev/bin/" OR "chown root:root x;" OR "/bin/telnet locip locport < /dev/console | /bin/sh" OR "/tmp/ratload" OR "ewok -t " OR "xspy -display " OR "cat > /dev/tcp/127.0.0.1/80 <<END" OR "rm -f /current/tmp/ftshell.latest" OR "ghost_* -v " OR " --wipe > /dev/null" OR "ping -c 2 *; grep * /proc/net/arp >/tmp/gx" OR "iptables * OUTPUT -p tcp -d 127.0.0.1 --tcp-flags RST RST -j DROP;" OR "> /var/log/audit/audit.log; rm -f ." OR "cp /var/log/audit/audit.log .tmp" OR "sh >/dev/tcp/* <&1 2>&1" OR "ncat -vv -l -p * <" OR "nc -vv -l -p * <" OR "< /dev/console | uudecode && uncompress" OR "sendmail -osendmail;chmod +x sendmail" OR "/usr/bin/wget -O /tmp/a http* && chmod 755 /tmp/cron" OR "chmod 666 /var/run/utmp~" OR "chmod 700 nscd crond" OR "cp /etc/shadow /tmp/." OR "</dev/console |uudecode > /dev/null 2>&1 && uncompress" OR "chmod 700 jp&&netstat -an|grep" OR "uudecode > /dev/null 2>&1 && uncompress -f * && chmod 755" OR "chmod 700 crond" OR "wget http*; chmod +x /tmp/sendmail" OR "chmod 700 fp sendmail pt" OR "chmod 755 /usr/vmsys/bin/pipe" OR "chmod -R 755 /usr/vmsys" OR "chmod 755 $opbin/*tunnel" OR "chmod 700 sendmail" OR "chmod 0700 sendmail" OR "/usr/bin/wget http*sendmail;chmod +x sendmail;" OR "&& telnet * 2>&1 </dev/console" | collect index=notable_events source="Equation Group Indicators" marker="guid=41e5c73d-9983-4b69-bd03-e13b67e9623c,tags=attack.execution,tags=attack.g0020,tags=attack.t1059.004,"
[Suspicious Use of /dev/tcp]
description = Detects suspicious command with /dev/tcp
search = index=syslog _index_earliest=-1h@h "cat </dev/tcp/" OR "exec 3<>/dev/tcp/" OR "echo >/dev/tcp/" OR "bash -i >& /dev/tcp/" OR "sh -i >& /dev/udp/" OR "0<&196;exec 196<>/dev/tcp/" OR "exec 5<>/dev/tcp/" OR "(sh)0>/dev/tcp/" OR "bash -c 'bash -i >& /dev/tcp/" OR "echo -e '#!/bin/bash\\nbash -i >& /dev/tcp/" | collect index=notable_events source="Suspicious Use of /dev/tcp" marker="guid=6cc5fceb-9a71-4c23-aeeb-963abe0b279c,tags=attack.reconnaissance,"
[Buffer Overflow Attempts]
description = Detects buffer overflow attempts in Unix system log files
search = index=syslog _index_earliest=-1h@h "attempt to execute code on stack by" OR "FTP LOGIN FROM .* 0bin0sh" OR "rpc.statd[\\d+]: gethostbyname error for" OR "AAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAAA" | collect index=notable_events source="Buffer Overflow Attempts" marker="guid=18b042f0-2ecd-4b6e-9f8d-aa7a7e7de781,tags=attack.t1068,tags=attack.privilege-escalation,"
[Commands to Clear or Remove the Syslog - Builtin]
description = Detects specific commands commonly used to remove or empty the syslog
search = index=syslog _index_earliest=-1h@h "rm /var/log/syslog" OR "rm -r /var/log/syslog" OR "rm -f /var/log/syslog" OR "rm -rf /var/log/syslog" OR "mv /var/log/syslog" OR " >/var/log/syslog" OR " > /var/log/syslog" NOT "/syslog." | collect index=notable_events source="Commands to Clear or Remove the Syslog - Builtin" marker="guid=e09eb557-96d2-4de9-ba2d-30f712a5afd3,tags=attack.impact,tags=attack.t1565.001,"
[Space After Filename]
description = Detects space after filename
search = index=syslog _index_earliest=-1h@h "echo \"*\" > * && chmod +x *" "mv * \"* \"" | collect index=notable_events source="Space After Filename" marker="guid=879c3015-c88b-4782-93d7-07adf92dbcb7,tags=attack.execution,"
[Suspicious Activity in Shell Commands]
description = Detects suspicious shell commands used in various exploit codes (see references)
search = index=syslog _index_earliest=-1h@h "wget * - http* | perl" OR "wget * - http* | sh" OR "wget * - http* | bash" OR "python -m SimpleHTTPServer" OR "-m http.server" OR "import pty; pty.spawn*" OR "socat exec:*" OR "socat -O /tmp/*" OR "socat tcp-connect*" OR "*echo binary >>*" OR "*wget *; chmod +x*" OR "*wget *; chmod 777 *" OR "*cd /tmp || cd /var/run || cd /mnt*" OR "*stop;service iptables stop;*" OR "*stop;SuSEfirewall2 stop;*" OR "chmod 777 2020*" OR "*>>/etc/rc.local" OR "*base64 -d /tmp/*" OR "* | base64 -d *" OR "*/chmod u+s *" OR "*chmod +s /tmp/*" OR "*chmod u+s /tmp/*" OR "* /tmp/haxhax*" OR "* /tmp/ns_sploit*" OR "nc -l -p *" OR "cp /bin/ksh *" OR "cp /bin/sh *" OR "* /tmp/*.b64 *" OR "*/tmp/ysocereal.jar*" OR "*/tmp/x *" OR "*; chmod +x /tmp/*" OR "*;chmod +x /tmp/*" | collect index=notable_events source="Suspicious Activity in Shell Commands" marker="guid=2aa1440c-9ae9-4d92-84a7-a9e5f5e31695,tags=attack.execution,tags=attack.t1059.004,"
[Linux Command History Tampering]
description = Detects commands that try to clear or tamper with the Linux command history. This technique is used by threat actors in order to evade defenses and execute commands without them being recorded in files such as "bash_history" or "zsh_history". 
search = index=syslog _index_earliest=-1h@h "cat /dev/null >*sh_history" OR "cat /dev/zero >*sh_history" OR "chattr +i*sh_history" OR "echo \"\" >*sh_history" OR "empty_bash_history" OR "export HISTFILESIZE=0" OR "history -c" OR "history -w" OR "ln -sf /dev/null *sh_history" OR "ln -sf /dev/zero *sh_history" OR "rm *sh_history" OR "shopt -ou history" OR "shopt -uo history" OR "shred *sh_history" OR "truncate -s0 *sh_history" | collect index=notable_events source="Linux Command History Tampering" marker="guid=fdc88d25-96fb-4b7c-9633-c0e417fdbd4e,tags=attack.defense-evasion,tags=attack.t1070.003,"
[Code Injection by ld.so Preload]
description = Detects the ld.so preload persistence file. See `man ld.so` for more information.
search = index=syslog _index_earliest=-1h@h "/etc/ld.so.preload" | collect index=notable_events source="Code Injection by ld.so Preload" marker="guid=7e3c4651-c347-40c4-b1d4-d48590fdf684,tags=attack.persistence,tags=attack.privilege-escalation,tags=attack.t1574.006,"
[JexBoss Command Sequence]
description = Detects suspicious command sequence that JexBoss
search = index=syslog _index_earliest=-1h@h "bash -c /bin/bash" "&/dev/tcp/" | collect index=notable_events source="JexBoss Command Sequence" marker="guid=8ec2c8b4-557a-4121-b87c-5dfb3a602fae,tags=attack.execution,tags=attack.t1059.004,"
[Suspicious Log Entries]
description = Detects suspicious log entries in Linux log files
search = index=syslog _index_earliest=-1h@h "entered promiscuous mode" OR "Deactivating service" OR "Oversized packet received from" OR "imuxsock begins to drop messages" | collect index=notable_events source="Suspicious Log Entries" marker="guid=f64b6e9a-5d9d-48a5-8289-e1dd2b3876e1,tags=attack.impact,"
[Remote File Copy]
description = Detects the use of tools that copy files from or to remote systems
search = index=syslog _index_earliest=-1h@h "scp " OR "rsync " OR "sftp " "@" OR ":" | collect index=notable_events source="Remote File Copy" marker="guid=7a14080d-a048-4de8-ae58-604ce58a795b,tags=attack.command-and-control,tags=attack.lateral-movement,tags=attack.t1105,"
[Guacamole Two Users Sharing Session Anomaly]
description = Detects suspicious session with two users present
search = index=syslog _index_earliest=-1h@h "(2 users now present)" | collect index=notable_events source="Guacamole Two Users Sharing Session Anomaly" marker="guid=1edd77db-0669-4fef-9598-165bda82826d,tags=attack.credential-access,tags=attack.t1212,"
[Modifying Crontab]
description = Detects suspicious modification of crontab file.
search = index=syslog _index_earliest=-1h@h "REPLACE" | collect index=notable_events source="Modifying Crontab" marker="guid=af202fd3-7bff-4212-a25a-fb34606cfcbe,tags=attack.persistence,tags=attack.t1053.003,"
[PwnKit Local Privilege Escalation]
description = Detects potential PwnKit exploitation CVE-2021-4034 in auth logs
search = index=syslog _index_earliest=-1h@h "pkexec" "The value for environment variable XAUTHORITY contains suspicious content" "[USER=root] [TTY=/dev/pts/0]" | collect index=notable_events source="PwnKit Local Privilege Escalation" marker="guid=0506a799-698b-43b4-85a1-ac4c84c720e9,tags=attack.privilege-escalation,tags=attack.t1548.001,"
[Suspicious VSFTPD Error Messages]
description = Detects suspicious VSFTPD error messages that indicate a fatal or suspicious error that could be caused by exploiting attempts
search = index=syslog _index_earliest=-1h@h "Connection refused: too many sessions for this address." OR "Connection refused: tcp_wrappers denial." OR "Bad HTTP verb." OR "port and pasv both active" OR "pasv and port both active" OR "Transfer done (but failed to open directory)." OR "Could not set file modification time." OR "bug: pid active in ptrace_sandbox_free" OR "PTRACE_SETOPTIONS failure" OR "weird status:" OR "couldn't handle sandbox event" OR "syscall * out of bounds" OR "syscall not permitted:" OR "syscall validate failed:" OR "Input line too long." OR "poor buffer accounting in str_netfd_alloc" OR "vsf_sysutil_read_loop" | collect index=notable_events source="Suspicious VSFTPD Error Messages" marker="guid=377f33a1-4b36-4ee1-acee-1dbe4b43cfbe,tags=attack.initial-access,tags=attack.t1190,"
[Relevant ClamAV Message]
description = Detects relevant ClamAV messages
search = index=syslog _index_earliest=-1h@h "Trojan*FOUND" OR "VirTool*FOUND" OR "Webshell*FOUND" OR "Rootkit*FOUND" OR "Htran*FOUND" | collect index=notable_events source="Relevant ClamAV Message" marker="guid=36aa86ca-fd9d-4456-814e-d3b1b8e1e0bb,tags=attack.resource-development,tags=attack.t1588.001,"
[Disabling Security Tools - Builtin]
description = Detects disabling security tools
search = index=syslog _index_earliest=-1h@h "stopping iptables" OR "stopping ip6tables" OR "stopping firewalld" OR "stopping cbdaemon" OR "stopping falcon-sensor" | collect index=notable_events source="Disabling Security Tools - Builtin" marker="guid=49f5dfc1-f92e-4d34-96fa-feba3f6acf36,tags=attack.defense-evasion,tags=attack.t1562.004,"
[Suspicious Named Error]
description = Detects suspicious DNS error messages that indicate a fatal or suspicious error that could be caused by exploiting attempts
search = index=syslog _index_earliest=-1h@h " dropping source port zero packet from " OR " denied AXFR from " OR " exiting (due to fatal error)" | collect index=notable_events source="Suspicious Named Error" marker="guid=c8e35e96-19ce-4f16-aeb6-fd5588dc5365,tags=attack.initial-access,tags=attack.t1190,"
[Sudo Privilege Escalation CVE-2019-14287 - Builtin]
description = Detects users trying to exploit sudo vulnerability reported in CVE-2019-14287
search = index=syslog _index_earliest=-1h@h USER IN ("#-*", "#*4294967295") | collect index=notable_events source="Sudo Privilege Escalation CVE-2019-14287 - Builtin" marker="guid=7fcc54cb-f27d-4684-84b7-436af096f858,tags=attack.privilege-escalation,tags=attack.t1068,tags=attack.t1548.003,tags=cve.2019-14287,"
[Suspicious OpenSSH Daemon Error]
description = Detects suspicious SSH / SSHD error messages that indicate a fatal or suspicious error that could be caused by exploiting attempts
search = index=syslog _index_earliest=-1h@h "unexpected internal error" OR "unknown or unsupported key type" OR "invalid certificate signing key" OR "invalid elliptic curve value" OR "incorrect signature" OR "error in libcrypto" OR "unexpected bytes remain after decoding" OR "fatal: buffer_get_string: bad string" OR "Local: crc32 compensation attack" OR "bad client public DH value" OR "Corrupted MAC on input" | collect index=notable_events source="Suspicious OpenSSH Daemon Error" marker="guid=e76b413a-83d0-4b94-8e4c-85db4a5b8bdc,tags=attack.initial-access,tags=attack.t1190,"
[SSHD Error Message CVE-2018-15473]
description = Detects exploitation attempt using public exploit code for CVE-2018-15473
search = index=syslog _index_earliest=-1h@h "error: buffer_get_ret: trying to get more bytes 1907 than in buffer 308 [preauth]" | collect index=notable_events source="SSHD Error Message CVE-2018-15473" marker="guid=4c9d903d-4939-4094-ade0-3cb748f4d7da,tags=attack.reconnaissance,tags=attack.t1589,"
